# dockerfiles/php/supervisor.conf
[supervisord]
nodaemon=true
loglevel = info
logfile=/home/logs/supervisor-logs/supervisord.log
pidfile=/var/run/supervisord.pid

[group:laravel-worker]
priority=999
programs=php8-fpm,schedule-run,laravel-queue

[program:php8-fpm]
priority=5
autostart=true
autorestart=true
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
command=/usr/local/sbin/php-fpm -R
stderr_logfile=/home/logs/supervisor-logs/php-error.log
stdout_logfile=/home/logs/supervisor-logs/php-access.log

[program:laravel-queue]
numprocs=5
autostart=true
autorestart=true
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/home/logs/supervisor-logs/worker.log
command=php /var/www/artisan queue:work --sleep=3 --tries=3

[program:schedule-run]
process_name=%(program_name)s_%(process_num)02d
command=/bin/bash -c "while [ true ]; do (php /var/www/artisan schedule:run --verbose --no-interaction && php /var/www/artisan short-schedule:run); sleep 60; done"
#command=/bin/bash -c "* * * * * cd /var/www/ && php artisan schedule:run >> /dev/null 2>&1"
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/home/logs/supervisor-logs/schedule.log
stopwaitsecs=60
